# deployment/supervisord.conf
; Process manager for the single Railway/Docker container: runs nginx, BOTH
; php-fpm pools (portfolio + jailed hosted), and the queue worker that executes
; ProcessHostedProject deploys. All logs go to stdout/stderr for the platform.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:php-fpm]
command=php-fpm -F
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g 'daemon off;'
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Runs the hosting deploy pipeline (security scan → extract → DB → nginx).
; --queue order makes hosting deploys take priority over default jobs.
[program:queue]
command=php /var/www/html/artisan queue:work --queue=hosting,default --tries=1 --timeout=600 --sleep=3 --max-time=3600
autorestart=true
priority=30
stopwaitsecs=610
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:scheduler]
command=/bin/sh -c "while true; do php /var/www/html/artisan schedule:run --no-interaction >> /dev/stdout 2>&1; sleep 60; done"
autorestart=true
priority=40
